# multi-market-bot.service
[Unit]
Description=Multi-Market Trading Bot v1.0
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/multi-market-bot
EnvironmentFile=/root/multi-market-bot/.env
ExecStart=/usr/bin/python3 /root/multi-market-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全退出
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
